# deb_package_files/DEBIAN/prerm
#!/bin/bash
systemctl disable KootnetEthServer
systemctl stop KootnetEthServer
crontab -u root -r
find /opt/kootnet-network-testers -type d -name __pycache__ -exec rm -rf {} \;
exit 0
